# color-scheme-watcher.service.in
[Unit]
Description=color-scheme-watcher.sh daemon
Documentation=https://github.com/alexbradd/color-scheme-watcher

[Service]
Type=simple
ExecStart=/bin/bash @PATH_TO_SCRIPT@/color-scheme-watcher.sh --verbose

[Install]
WantedBy=graphical-session.target
